# roles/gentooinstallbattle._common/templates/etc/profile.d/makeopts.sh.j2
#!/bin/bash
export MAKEOPTS="-j$(($(nproc)*2+1)) -l$(($(nproc)+1))"
