# templates/etc/profile.d/makeopts.sh.j2
#!/bin/bash -eu
export MAKEOPTS="-j$(($(nproc)+1))"
